# etc/init/pir-sensor.conf
description "PIR sensor"
author "Ken Chou <kenchou77@gmail.com>"

start on local-filesystems
stop on shutdown

respawn

exec /home/pi/pir-video-capture/bin/pir.py

