# etc/init/pir-sensor.conf
description "PIR sensor"
author "Ken Chou <kenchou77@gmail.com>"

start on startup
stop on shutdown

respawn

exec /home/pi/pir-video-capture/pir/pir.py

# post-stop exec stop video-capture
